# scripts/codicil-mcp.service.template
[Unit]
Description=Codicil MCP HTTP Server
After=network.target

[Service]
Type=simple
User=codicil
Group=codicil
WorkingDirectory=/opt/codicil
EnvironmentFile=/etc/codicil/codicil.env
ExecStart=/usr/bin/node /opt/codicil/scripts/mcp-server-http.mjs
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/codicil

[Install]
WantedBy=multi-user.target

