# scripts/memex-mcp.service.template
[Unit]
Description=Memex MCP HTTP Server
After=network.target

[Service]
Type=simple
User=memex
Group=memex
WorkingDirectory=/opt/memex
EnvironmentFile=/etc/memex/memex.env
ExecStart=/usr/bin/node /opt/memex/scripts/mcp-server-http.mjs
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/memex

[Install]
WantedBy=multi-user.target

